# htdocs/local/php_interface/crontab/crontab.conf
#! /usr/bin/php
*/10 * * * * /usr/bin/php -f /home/www/klavazip.ru/www/bitrix/php_interface/crontab/__crontab.php BX_CRONTAB=1 > /home/www/klavazip.ru/www/bitrix/php_interface/crontab/logs/cron.1.log
*/10 * * * * /usr/bin/php -f /home/www/klavazip.ru/www/bitrix/php_interface/crontab/__crontab.php > /home/www/klavazip.ru/www/bitrix/php_interface/crontab/logs/cron.2.log
*/15 * * * * /usr/bin/php -f /home/www/klavazip.ru/www/bitrix/php_interface/crontab/__crontab.php BX_SUBSCRIBE=1 > /home/www/klavazip.ru/www/bitrix/php_interface/crontab/logs/cron.3.log